# etc/supervisord.conf
[unix_http_server]
file=./var/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=./var/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=./var/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://./var/supervisor.sock ; use a unix:// URL  for a unix socket

[program:Mail]
command=./MailWorker
process_name=%(process_num)s
directory=.
stdout_logfile=./var/Mail.out
stderr_logfile=./var/Mail.err
numprocs=5
numprocs_start=5
startsecs=0
autorestart=true

[program:ImageUploader]
command=./ImageUploaderWorker
process_name=%(process_num)s
directory=.
stdout_logfile=./var/ImageUploader.out
stderr_logfile=./var/ImageUploader.err
numprocs=5
numprocs_start=5
startsecs=0
autorestart=true
